# com.gustave.reddit-watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.gustave.reddit-watcher</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/andreidonko/Gustave/gustave/reddit_watcher/run_scheduled.sh</string>
    </array>

    <!-- Run every 3 hours (10800 seconds). launchd skips a tick if the
         previous run is still going, so scans never overlap. -->
    <key>StartInterval</key>
    <integer>10800</integer>

    <!-- Also run once immediately when the job is first loaded, so you can
         confirm it works without waiting 3 hours. -->
    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/andreidonko/Gustave/gustave/reddit_watcher/reports/launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/andreidonko/Gustave/gustave/reddit_watcher/reports/launchd.err.log</string>
</dict>
</plist>
